# com.google-docs-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.google-docs-mcp</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/afridshaik/google-docs-mcp/server-http.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/afridshaik/google-docs-mcp</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/google-docs-mcp.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/google-docs-mcp.err</string>
</dict>
</plist>
